# docker/crontab
# START CRON JOB LIST
* 8-22 * * * python3 /var/tmp/src/KSL-Items/getItemList.py > /var/cron.log
# END CRON JOB LIST