# docker/crontab
# START CRON JOB LIST
* 8-22 * * * python3 /var/tmp/src/KSL-Items/getItemListDocker.py > /var/cron.log
#* 8-22 * * * python3 /var/tmp/src/KSL-Cars/getCarDocker.py > /var/cron.log
# END CRON JOB LIST